# rocon/devel/etc/catkin/profile.d/15.rocon_test.bash
# generated from rocon_test/env-hooks/15.rocon_test.bash.em

. "/home/chris/Documents/mines/mines_grad/csci598/hcr_term_project/rocon/src/rocon_multimaster/rocon_test/shells/env.bash"
