# rocon/devel/etc/catkin/profile.d/15.rocon_app.tcsh
# generated from rocon_app_utilities/env-hooks/15.rocon_app_utilities.tcsh.em

. "/home/chris/Documents/mines/mines_grad/csci598/hcr_term_project/rocon/src/rocon_app_platform/rocon_app_utilities/shells/env.tcsh"
