# rocon/devel/etc/catkin/profile.d/15.rospair.bash
# generated from rocon_python_comms/env-hooks/15.rospair.bash.em

. "/home/chris/Documents/mines/mines_grad/csci598/hcr_term_project/rocon/src/rocon_tools/rocon_python_comms/completion/rospairbash"
