# turtlebot/devel/etc/catkin/profile.d/15.ld_lib_path.bash
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/home/marshall/hcr_term_project/turtlebot/install/lib/python2.7/dist-packages
